# distr/echoBot.service
[Unit]
Description=Echo Bot Service
After=multi-user.target
After=syslog.target
After=network.target
Conflicts=getty@tty1.service

[Service]
Type=simple
User=kulyashovv
WorkingDirectory=/home/kulyashovv/MLsas/
ExecStart=/usr/bin/python3 /home/kulyashovv/MLsas/echo_bot.py
StandardInput=tty-force
Restart=on-failure

[Install]
WantedBy=multi-user.target
